# php/crontabs/www-data
* * * * * /usr/local/bin/php /var/www/html/nb-api/artisan schedule:run >> /dev/null 2>&1
